# deploy/btcedu-web.service
[Unit]
Description=btcedu web dashboard (gunicorn)
After=network.target

[Service]
Type=notify
User=pi
WorkingDirectory=/home/pi/AI-Startup-Lab/bitcoin-education
ExecStart=/home/pi/AI-Startup-Lab/bitcoin-education/.venv/bin/gunicorn \
    -w 2 \
    -b 127.0.0.1:8091 \
    --timeout 300 \
    "btcedu.web.app:create_app()"
EnvironmentFile=/home/pi/AI-Startup-Lab/bitcoin-education/.env
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
